# etc/StoRM_watchdog.cron
* * * * * /opt/storm/etc/StoRM_watchdog.sh
