# etc/cron.d/storm-backend-server.cron
00 00 * * * root /usr/sbin/logrotate -s /etc/logrotate.d/logrotate.status /etc/logrotate.d/storm-backend-server
