# etc/cron.d/storm-backend.cron
00 00 * * * root /usr/sbin/logrotate -s /opt/storm/backend/etc/logrotate.d/logrotate.status /opt/storm/backend/etc/logrotate.d/storm-backend.logrotate
